# contrib/buddycloud-media-server.conf
# buddycloud-http-api - buddycloud Media Server

description     "buddycloud Media Server"

# respawn the job up to 10 times within a 5 second period.
# If the job exceeds these values, it will be stopped and
# marked as failed.
respawn
respawn limit 10 5

start on runlevel [2345]
stop on runlevel [!2345]

setuid nobody
setgid nogroup

chdir /opt/buddycloud-media-server
exec java -cp .:*:lib/* com.buddycloud.mediaserver.Main
